# server/service/mysql/run
#!/bin/bash

exec 2>&1
exec mysqld --user=mysql --log-error=/var/log/mysql/error.log --pid-file=/var/lib/mysql/mysqld.pid --port=3306
