# server/service/graphite_exporter/finish
#!/bin/bash

/usr/bin/s6-svscanctl -t /service
