# server/service/graphite_exporter/run
#!/bin/bash

exec 2>&1
exec /usr/bin/graphite_exporter
