# server/service/cattle/run
#!/bin/bash

exec /usr/share/cattle/cattle.sh
