# cmd/rancherd/installer/etc/profile.d/rke2.sh
KUBECONFIG="/etc/rancher/rke2/rke2.yaml"
PATH="$PATH:/var/lib/rancher/rke2/bin"
