# server/service/pgsql/run
#!/bin/bash

su-exec postgres /usr/bin/postgres -D /var/lib/pgsql/9.6
